# kasa_mqtt/bin/kasa-mqtt.service
[Unit]
Description=MQTT front end wrapper to python-kasa

[Service]
User=pi
Type=simple
ExecStart=/home/pi/kasa-mqtt/kasa_mqtt/bin/start_kasa-mqtt.sh /home/pi/kasa-mqtt/config/config.yaml
Restart=on-failure

[Install]
WantedBy=multi-user.target
